# packaging/pewpawd.service
[Unit]
Description=PewPaw WhatsApp Backend Daemon
Documentation=https://github.com/COXKPER/PewPaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pewpawd
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
